# systemd/vps-probe.service
[Unit]
Description=Jager Monitor Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vps-probe
ExecStart=/usr/bin/python3 /opt/vps-probe/panel.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
MemoryMax=128M
TasksMax=128
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
